# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=/usr/sbin/php-fpm8 -F
autostart=true
autorestart=true
priority=999
stdout_logfile=/var/log/php-fpm.log
stderr_logfile=/var/log/php-fpm-error.log
user=nobody

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=998
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx-error.log

[program:laravel-scheduler]
process_name=%(program_name)s
command=php /app/artisan schedule:run
autostart=true
autorestart=true
priority=997
stdout_logfile=/var/log/scheduler.log
stderr_logfile=/var/log/scheduler-error.log
startsecs=0
stopwaitsecs=3600
redirect_stderr=true
user=nobody

[program:laravel-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /app/artisan queue:work --tries=3 --timeout=300
autostart=true
autorestart=true
numprocs=2
priority=996
stdout_logfile=/var/log/queue-%(process_num)02d.log
stderr_logfile=/var/log/queue-error-%(process_num)02d.log
redirect_stderr=true
user=nobody
